# config/movie.conf
[program:movie]
command=/home/john/venv/bin/gunicorn TEB.wsgi:application -c /home/john/django-telegram-auth-example/config/gunicorn.conf.py
directory=/home/john/django-telegram-auth-example
user=john
autorestart=true
redirect_stderr=true
stdout_logfile=/home/john/django-telegram-auth-example/logs/debug.log
